# linux/systemd/.config/systemd/user/fish-update-completions.timer
[Unit]
Description=Run fish-update-completions daily and on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
