# autostart/.config/autostart/polybar.desktop
[Desktop Entry]
Exec=/home/gary/.config/polybar/launch.sh
Name=Launch Polybar
Type=Application
OnlyShowIn=i3
